# linear-chief.service
[Unit]
Description=Linear Chief of Staff - AI-powered briefing agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=%i
Group=%i
WorkingDirectory=/home/%i/linear-agent

# Environment
Environment="PATH=/home/%i/linear-agent/.venv/bin:/usr/local/bin:/usr/bin:/bin"
Environment="PYTHONUNBUFFERED=1"
EnvironmentFile=/home/%i/linear-agent/.env

# Execute scheduler
ExecStart=/home/%i/linear-agent/.venv/bin/python -m linear_chief start

# Restart policy
Restart=always
RestartSec=10
StartLimitInterval=200
StartLimitBurst=5

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=linear-chief

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/%i/.linear_chief
ReadWritePaths=/home/%i/linear-agent

# Resource limits
MemoryLimit=512M
CPUQuota=50%

[Install]
WantedBy=multi-user.target
